# overlay/etc/systemd/system/mimic-node.path
[Unit]
Description=Watch for Mimic-Node Staging Configuration Changes

[Path]
# Watch for changes to the staging config file.
# When mimictl creates/updates this file (e.g. via batch add),
PathChanged=/etc/sing-box/config.json
Unit=mimic-node-deploy.service

[Install]
WantedBy=multi-user.target
