# overlay/etc/systemd/system/mimic-node.service
[Unit]
Description=Mimic-Node Auto Maintenance (SNI Rotation & Config Apply)
Requires=mimic-node-mount.service
After=network-online.target mimic-node-mount.service
Wants=network-online.target

[Service]
Type=oneshot
# 1. Detect optimal SNI and apply any staged config (config.new)
#    This consumes config.new if it exists, merging it with new SNI.
ExecStart=/usr/bin/mimictl sni

[Install]
WantedBy=multi-user.target
